# com.leaderboard.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leaderboard.startup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPLACE_WITH_ABSOLUTE_PATH/start_local.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/leaderboard.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/leaderboard.out</string>
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_ABSOLUTE_PATH</string>
</dict>
</plist> 